# automation/launchd/net.jessetek.hourly-autopush.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>net.jessetek.hourly-autopush</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/jesseonate/Documents/Claude Code/Jessetek/landing-page/automation/scripts/run-hourly-autopush.sh</string>
  </array>
  <key>StartInterval</key><integer>3600</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/jesseonate/Documents/Claude Code/Jessetek/landing-page/automation/logs/hourly-autopush.out</string>
  <key>StandardErrorPath</key><string>/Users/jesseonate/Documents/Claude Code/Jessetek/landing-page/automation/logs/hourly-autopush.err</string>
</dict></plist>
